# deploy/s6/services.d/kubecode/run
#!/command/with-contenv bash
set -euo pipefail

/etc/init-persistent-home
cd "${PERSISTENT_DIR}"

export NB_PREFIX="${NB_PREFIX:-/}"
export PORT=8888
export HOST=0.0.0.0
export KUBECODE_STATIC_DIR="${KUBECODE_STATIC_DIR:-/opt/kubecode/dist}"

echo "INFO: starting Kubecode at ${NB_PREFIX}"
exec 2>&1
exec /usr/local/bin/kubecode-server
